# alpine-s6-nginx-uwsgi/etc/s6-overlay/s6-rc.d/uwsgi-init/up
with-contenv

importas -D "" UWSGI_LOG UWSGI_LOG

ifelse -Xn { test -n "${UWSGI_LOG}" } { exit }

case ${UWSGI_LOG} {
	"yes|true|on" {
		foreground { echo "uwsgi-init: info: UWSGI_LOG enabled" }
		sed -e "s|^disable-logging.*|disable-logging = false|" -i "/etc/uwsgi.ini"
	}
}
echo "uwsgi-init: info: UWSGI_LOG disabled"
