# debian-slim-s6-nginx-php-fpm/root/etc/s6-overlay/s6-rc.d/php-fpm/run
#! /usr/bin/execlineb -P

with-contenv
importas -D "php-fpm7" php_binary PHP_BINARY
${php_binary}
