# alpine-s6-nginx-php-fpm/etc/s6-overlay/s6-rc.d/php-fpm-init/up
with-contenv

#	configure timezone
#
foreground
{
	importas -D "" TZ TZ
	ifelse -X { test -n "${TZ}" }
	{
		importas -D "" PHP_PACKAGE PHP_PACKAGE
		define php_ini "/etc/${PHP_PACKAGE}/php.ini"

		ifelse -Xn { test -f "${php_ini}" }
		{
			foreground { echo "php-fpm-init: ERROR: can't find PHP config: ${php_ini}" }
			exit 1
		}
		foreground { sed -E -i "${php_ini}" -e "s|^;?\\s*date\\.timezone.*|date.timezone = ${TZ}|g" }
		echo "php-fpm-init: info: set timezone: ${TZ}"
	}
}
