# debian-slim-s6/root/etc/s6-overlay/s6-rc.d/timezone/up
with-contenv
importas -D "" TZ TZ
ifelse -Xn { eltest -n "${TZ}" } { exit }

define tz_file "/usr/share/zoneinfo/${TZ}"

ifelse -X { eltest -f "${tz_file}" }
  {
    foreground { ln -sf "${tz_file}" "/etc/localtime" }
    echo "timezone: info: set timezone: ${TZ}"
  }
  foreground { echo "timezone: ERROR: cannot find ${tz_file}" }
  echo "timezone: ERROR: invalid TZ string: ${TZ}"
